# launchd/io.etherwise.v3.upwork-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.etherwise.v3.upwork-scan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/abhijeet/Desktop/Etherwise/etherwise-v3 &amp;&amp; PYTHONPATH=. /usr/bin/python3 -m modules.upwork.scan_pipeline >> var/logs/launchd-upwork-scan.log 2>&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
