# launchd/io.etherwise.v3.upwork-outcomes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.etherwise.v3.upwork-outcomes</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Developer/CommandLineTools/Library/Frameworks/Python3.framework/Versions/3.9/Resources/Python.app/Contents/MacOS/Python</string>
    <string>-m</string>
    <string>modules.upwork.outcome_capture</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/abhijeet/Desktop/Etherwise/etherwise-v3</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/abhijeet/Desktop/Etherwise/etherwise-v3</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/abhijeet/Desktop/Etherwise/etherwise-v3/var/logs/launchd-upwork-outcomes.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/abhijeet/Desktop/Etherwise/etherwise-v3/var/logs/launchd-upwork-outcomes.log</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>11</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
