# launchd/io.etherwise.v3.cockpit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>io.etherwise.v3.cockpit</string>
  <!-- The one always-on localhost server the architecture allows. Framework
       Python (FDA), NO WorkingDirectory, logs OUTSIDE ~/Desktop (June-14 TCC
       pattern). KeepAlive = relaunch if it ever exits. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Developer/CommandLineTools/Library/Frameworks/Python3.framework/Versions/3.9/Resources/Python.app/Contents/MacOS/Python</string>
    <string>-m</string>
    <string>modules.cockpit.server</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>PYTHONPATH</key><string>/Users/abhijeet/Desktop/Etherwise/etherwise-v3</string>
  </dict>
  <key>StandardOutPath</key><string>/Users/abhijeet/Library/Logs/Etherwise-v3/cockpit.log</string>
  <key>StandardErrorPath</key><string>/Users/abhijeet/Library/Logs/Etherwise-v3/cockpit.log</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
</dict>
</plist>
